# Linux/Misc/remove_extra_firm.service
[Unit]
Description=RemoveExtra

[Service]
Type=forking
ExecStart=/usr/bin/remove_extra_firm.sh

[Install]
WantedBy=graphical.target
